# iatrade.service
# Systemd Service File - Para Oracle Free Tier
# Copie este arquivo para /etc/systemd/system/iatrade.service na sua instância Oracle

# Como usar:
# 1. Fazer SSH na Oracle
# 2. sudo nano /etc/systemd/system/iatrade.service
# 3. Colar o conteúdo abaixo
# 4. Ctrl+X, Y, Enter
# 5. sudo systemctl daemon-reload
# 6. sudo systemctl enable iatrade.service
# 7. sudo systemctl start iatrade.service
# 8. sudo systemctl status iatrade.service

[Unit]
Description=IaTrade Bot - Trading Automático 24h
After=network.target
Wants=network-online.target

[Service]
# Identidade do serviço
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/IaTrade

# Ambiente
Environment="PATH=/home/ubuntu/IaTrade/venv/bin"
Environment="PYTHONUNBUFFERED=1"
EnvironmentFile=/home/ubuntu/IaTrade/.env

# Comando a executar
ExecStart=/home/ubuntu/IaTrade/venv/bin/python start_bot_24h.py

# Comportamento em caso de falha
Restart=always
RestartSec=60
StartLimitInterval=3600
StartLimitBurst=5

# Logs
StandardOutput=journal
StandardError=journal
SyslogIdentifier=iatrade

# Limites
LimitNOFILE=65536
LimitNPROC=4096

# Timeout gracioso
TimeoutStopSec=30
KillMode=mixed

[Install]
WantedBy=multi-user.target
